# ha-raumkernel-addon/rootfs/etc/services.d/raumkernel/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Raumkernel service
# ==============================================================================

bashio::log.info "Starting Raumkernel..."

## Run your program
cd /app
exec node index.js
